# mousemode.service
[Unit]
Description=mousemode

[Service]
ExecStart=/home/adrien/travail/mousemode/mousemode

[Install]
WantedBy=multi-user.target
